# raspberrypi/rc.local

# This commands need to be added to rc.local@pi


# Create a FIFO to receive commands and send to CEC-CLIENT
/home/pi/scripts/start.sh

# switch TV input to cable box (if someone turn on the TV, the input will be at cable box).
/home/pi/scripts/fastnet.sh

# standby the TV...
/home/pi/scripts/standby.sh

# Print the IP address

_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

echo
echo
echo
echo ' Bom dia Lucas! Aproveite a musica!';
echo
echo
echo

